# deploy/pi/systemd/promtail-log-access.service
[Unit]
Description=Allow non-root Promtail to enumerate Docker container logs

[Service]
Type=oneshot
ExecStart=/opt/song-history/scripts/prepare-promtail-log-access.sh
