# deploy/pi/systemd/promtail-log-access.path
[Unit]
Description=Watch for Docker container directories needed by Promtail

[Path]
PathChanged=/var/lib/docker/containers
Unit=promtail-log-access.service

[Install]
WantedBy=multi-user.target
